# templates/wol.service.j2
[Unit]
Description=Enable Wake On Lan on %i

[Service]
Type=oneshot
ExecStart=/usr/sbin/ethtool --change %i wol g
RemainAfterExit=true
StandardOutput=journal

[Install]
WantedBy=multi-user.target
