# templates/wol.service.j2
[Unit]
Description=Enable Wake On Lan on {{ item }}

[Service]
Type=oneshot
ExecStart = <Your filepath> --change {{ item }} wol g

[Install]
WantedBy=basic.target
